# config/link/waybar/waybar.service
[Unit]
Description=Start/Stop Waybar

[Service]
ExecStart=/home/diverseuniverse/startWaybar
Restart=on-failure

[Install]
WantedBy=default.target
